# Chapter14/etc/sv/haveged/run
#!/bin/sh
exec /usr/sbin/haveged -w 1024 -r 0 -F
