# Chapter13/simpleserver-systemd/simpleserver.service
[Unit] 
Description=Simple server 
 
[Service] 
Type=forking 
ExecStart=/usr/bin/simpleserver 
 
[Install] 
WantedBy=multi-user.target

